# roles/backup-webserver/files/restic_backup.cron
#!/bin/sh

#
# Cronjob to perform system backup
#


# run with low priority:
ionice -n 7 -p $$
renice 3 -p $$ > /dev/null


# create backup of InfluxDB:
rm -rf /var/backups/influxdb/
influxd backup -portable /var/backups/influxdb/ >/dev/null
# unpack backup, for better deduplication by Restic:
for f in /var/backups/influxdb/*.tar.gz; do
	gunzip "$f"
done


# TODO: also create backup of MySQL data


# do the actual restic invocations:
run-parts /etc/backup.d/
